# docker-jans-all-in-one/app/conf/supervisord.conf
[supervisord]
# logfile = /dev/stdout
# logfile_maxbytes = 0
nodaemon = true

[supervisorctl]

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command = nginx -g 'daemon off; pid /run/nginx/nginx.pid;'
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0
startretries = 10
priority = 1000

[include]
files = {supervisord_conf_includes}
