# docker-jans-all-in-one/app/conf/supervisord.conf.d/jans-saml.conf
[program:jans-saml]
command = sh /app/bin/jans-saml-entrypoint.sh
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/1
stderr_logfile_maxbytes = 0
startretries = 10
environment = CN_SAML_JAVA_OPTIONS="-Xms{min_heap_size}m -Xmx{max_heap_size}m {java_options}"
