# debian/oxd-https-extension
### BEGIN INIT INFO
# Provides:          oxd-https-extension
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

CONF=/opt/oxd-https-extension/lib/oxd-https.yml
OXD_HTTPS_EXTENSION_VERSION="4.0-SNAPSHOT"
SERVICE_NAME=oxd-https-extension
PID_PATH_NAME=/var/run/oxd-https-extension.pid
LIB=/opt/oxd-https-extension/lib/oxd-https-extension-${OXD_HTTPS_EXTENSION_VERSION}.jar


do_start () {
        if [ ! -f $PID_PATH_NAME ]; then
                echo "Starting $SERVICE_NAME ..."
                nohup java -jar $LIB server $CONF 2>>/dev/null>>/dev/null&
                echo $! > $PID_PATH_NAME
        else
            echo "$SERVICE_NAME is already running ..."
        fi
        PID_NUM=$(cat $PID_PATH_NAME)
        echo "PID: [$PID_NUM]"
}

do_stop () {
        if [ -f $PID_PATH_NAME ]; then
            PID=$(cat $PID_PATH_NAME);
            echo "$SERVICE_NAME stopping ..."
            kill $PID;
            rm $PID_PATH_NAME
        else
            echo "$SERVICE_NAME is not running ..."
        fi
}

case $1 in
    start)
        do_start
    ;;
    stop)
        do_stop
    ;;
    restart)
        do_stop
        do_start
    ;;
    status)
        if [ -f $PID_PATH_NAME ]; then
                echo "$SERVICE_NAME is running ...";
                PID_NUM=$(cat $PID_PATH_NAME)
                echo "PID: [$PID_NUM]"
        else
                echo "$SERVICE_NAME is not running ..."
        fi
    ;;
    *)
        echo "Usage: $0 {start|stop|status|restart}"
        RETVAL=2
    ;;
esac
