# debian/postinst
#!/bin/sh -e
update-rc.d gluu-oxd-server defaults > /dev/null
