# debian/postinst
#!/bin/sh -e
update-rc.d gluu-oxd-local defaults > /dev/null
