# jans-linux-setup/jans_setup/static/system/systemd/jans-link.service
[Unit]
Description=Janssen Link Server
After=%(order_jans_link_service)s

[Service]
Type=forking
ExecStart=/opt/dist/scripts/jans-link start
ExecStop=/opt/dist/scripts/jans-link stop
User=jetty
Group=jetty

[Install]
WantedBy=multi-user.target
